# supervisord.conf
[supervisord]
nodaemon=true
user=1000
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
logfile_maxbytes=50MB
logfile_backups=10

[program:reset-login]
command=/usr/bin/python3 /opt/nexus/force_setup_admin.py
directory=/opt/nexus
autorestart=false
startsecs=0
priority=1
stdout_logfile=/tmp/reset-login.log
stderr_logfile=/tmp/reset-login.err

[program:evolution]
command=/bin/bash /opt/nexus/start_evolution.sh
directory=/opt/nexus
autorestart=true
startretries=10
startsecs=5
stdout_logfile=/tmp/evolution.log
stderr_logfile=/tmp/evolution.err
environment=PORT="8080",DATABASE_ENABLED="true",NODE_ENV="production",MANAGER_ENABLED="true",WHATSAPP_IDE_ENABLED="true",AUTHENTICATION_TYPE="apikey",AUTHENTICATION_API_KEY="Supera",AUTHENTICATION_EXPOSE_IN_FETCH_INSTANCES="true",DATABASE_CONNECTION_URI="postgresql://postgres.htabdguydyysolkzdilm:*Mm0101mM****@aws-0-us-west-2.pooler.supabase.com:5432/postgres?schema=evolution_api",PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

[program:n8n]
command=n8n start
directory=/opt/nexus
priority=50
autorestart=true
stdout_logfile=/tmp/n8n.log
stderr_logfile=/tmp/n8n.err
environment=N8N_PORT="5678",N8N_LISTEN_ADDRESS="127.0.0.1",N8N_ENCRYPTION_KEY="nexus_secret_key",N8N_USER_FOLDER="/opt/nexus/.n8n",HOME="/opt/nexus",N8N_PROTOCOL="http",WEBHOOK_URL="https://info86-new-espacio-jazmio.hf.space/n8n/",N8N_BASE_URL="https://info86-new-espacio-jazmio.hf.space/n8n/",N8N_EDITOR_BASE_URL="https://info86-new-espacio-jazmio.hf.space/n8n/",N8N_LOG_LEVEL="debug",N8N_BLOCK_IFRAME_ACCESS="false",PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

[program:status]
command=/usr/bin/python3 -m uvicorn status:app --host 0.0.0.0 --port 8001
directory=/opt/nexus
priority=3
autorestart=true
stdout_logfile=/tmp/status.log
stderr_logfile=/tmp/status.err

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
priority=100
autorestart=true
stdout_logfile=/tmp/nginx.log
stderr_logfile=/tmp/nginx.err

[program:redis]
command=redis-server --protected-mode no
priority=1
autorestart=true
stdout_logfile=/tmp/redis.log
stderr_logfile=/tmp/redis.err

[program:keepalive]
command=/usr/bin/python3 /opt/nexus/keepalive.py
directory=/opt/nexus
autorestart=true
stdout_logfile=/tmp/keepalive.log
stderr_logfile=/tmp/keepalive.err
