# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
loglevel=info

; ============================================================
; STATUS DASHBOARD (Puerto 8001)
; ============================================================
[program:status]
command=python3 status.py
autostart=true
autorestart=true
startsecs=3
startretries=5
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PORT="8001",DATABASE_URL="postgresql://postgres.htabdguydyysolkzdilm:*Mm0101mM****@aws-0-us-west-2.pooler.supabase.com:5432/postgres?schema=evolution_api"

; ============================================================
; REDIS (Requerido por Evolution API)
; ============================================================
[program:redis]
command=/usr/bin/redis-server --protected-mode no --dir /opt/nexus/redis/
autostart=true
autorestart=true
startsecs=2
startretries=5
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ============================================================
; EVOLUTION API (Puerto 8080)
; ============================================================
[program:evolution]
command=/bin/bash start_evolution.sh
directory=/opt/nexus
autostart=true
autorestart=true
startsecs=10
startretries=10
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PORT="8080",SERVER_URL="https://jazm-io-production.up.railway.app",DEL_INSTANCE="false"

; ============================================================
; N8N AUTOMATION (Puerto 3100 Master)
; ============================================================
[program:n8n]
command=/bin/bash start_n8n.sh
directory=/opt/nexus
autostart=true
autorestart=true
startsecs=30
startretries=20
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=N8N_PORT="3100",N8N_PATH="/n8n/",N8N_PROTOCOL="http",N8N_EDITOR_BASE_URL="https://jazm-io-production.up.railway.app/n8n/",WEBHOOK_URL="https://jazm-io-production.up.railway.app/n8n/",N8N_USER_FOLDER="/opt/nexus/.n8n",EXECUTIONS_PROCESS="main",N8N_ENCRYPTION_KEY="nexus_secret_key",DB_TYPE="postgresdb",DB_POSTGRESDB_HOST="aws-0-us-west-2.pooler.supabase.com",DB_POSTGRESDB_PORT="5432",DB_POSTGRESDB_DATABASE="postgres",DB_POSTGRESDB_USER="postgres.htabdguydyysolkzdilm",DB_POSTGRESDB_PASSWORD="*Mm0101mM****",DB_POSTGRESDB_SCHEMA="n8n_fresh_access",DB_POSTGRESDB_SSL="true",DB_POSTGRESDB_SSL_REJECT_UNAUTHORIZED="false",N8N_ENFORCE_SETTINGS_FILE_PERMISSIONS="false"

; ============================================================
; NGINX (Puerto 7860 Gateway)
; ============================================================
[program:nginx]
command=/usr/sbin/nginx -g "daemon off; error_log /tmp/nginx/logs/error.log; pid /tmp/nginx/nginx.pid;"
autostart=true
autorestart=true
startsecs=5
startretries=5
priority=50
stdout_logfile=/tmp/nginx_out.log
stderr_logfile=/tmp/nginx_err.log
